# osp_sai_2.1.8/platform/build_svc/etc/hi3535/inittab
::sysinit:/etc/init.d/rcS
ttyAMA0::respawn:/usr/sbin/klish
#::respawn:/usr/sbin/klish
::shutdown:/usr/sbin/restart.sh