# osp_sai_2.1.8/platform/build_svc/etc/e350/inittab
::sysinit:/switch/etc/rc.sysinit
#ttyS0::respawn:/usr/sbin/klish
ttyS0::respawn:/switch/etc/start_klish_for_console.sh
#::respawn:/usr/sbin/klish
::shutdown:/usr/sbin/restart.sh
