# osp_sai_2.1.8/platform/build_svc/etc/e580/debian/debian_installer/80-mount-udisk.rules
ACTION!="add",GOTO="farsight"
KERNEL=="sd[a-z][0-9]",RUN+="/etc/mount.sh %k"
KERNEL=="mmcblk[0-9]p[0-9]",RUN+="/etc/mount.sh %k"
LABEL="farsight"
