# osp_sai_2.1.8/platform/build_svc/etc/udev/rules.d/unmount_x86_64.rules
ACTION!="remove",GOTO="farsight"
DEVNAME=="sd[a-z][0-9]",RUN+="/etc/unmount.sh"
LABEL="farsight"
