# osp_sai_2.1.8/platform/build_svc/etc/e582/inittab
::sysinit:/switch/etc/rc.sysinit
#console::respawn:-/bin/sh
ttyS0::respawn:/switch/etc/start_klish_for_console.sh
::shutdown:/usr/sbin/restart.sh
